# src/php-fpm/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * www-data php /var/www/htdocs/scripts/send_mails.php >> /var/www/logs/send_mails.log 2>&1
